# etc/init/localproxy.conf
# localproxy
description "node local proxy"
author "Tom O'Connor"
respawn 
start on runlevel [2345]

chdir /root
exec node localproxy.js >> /var/log/proxy.log 2>&1


