# etc/init/remoteproxy.conf
# remoteproxy
description "node remote proxy"
author "Tom O'Connor"
respawn 
start on runlevel [2345]

chdir /root
exec node remoteproxy.js >> /var/log/proxy.log 2>&1


